# install/mnttripbox
#!/bin/sh /etc/rc.common

START=50

start() {
	mount /dev/sda3 /usbstick
}
